# launchd/com.chat-daily-tg.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>REPLACE_WITH_HOME/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>Label</key>
	<string>com.chat-daily-tg.agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>REPLACE_WITH_PROJECT_DIR/bin/cdrun-bash</string>
		<string>REPLACE_WITH_PROJECT_DIR/scripts/run_daily_guarded.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>REPLACE_WITH_DATA_DIR/logs/stderr.log</string>
	<key>StandardOutPath</key>
	<string>REPLACE_WITH_DATA_DIR/logs/stdout.log</string>
	<key>StartCalendarInterval</key>
	<array>
		<!-- Single 07:05 trigger. Delivery timing is wake-gated in-process
		     (run_daily --wait-for-wake, 5-min polls, 13:00 deadline); a 07:05
		     slept through is coalesced by launchd and re-fired on Mac wake,
		     so the old 09:00/13:00 catch-up entries are gone. -->
		<dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>5</integer></dict>
	</array>
	<key>WorkingDirectory</key>
	<string>REPLACE_WITH_PROJECT_DIR</string>
</dict>
</plist>
