# launchd/com.chat-daily-tg.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/Apple/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
	</dict>
	<key>Label</key>
	<string>com.chat-daily-tg.agent</string>
	<key>ProgramArguments</key>
	<array>
		<string>REPLACE_WITH_PROJECT_DIR/.venv/bin/python</string>
		<string>REPLACE_WITH_PROJECT_DIR/run_daily.py</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>StandardErrorPath</key>
	<string>REPLACE_WITH_DATA_DIR/logs/stderr.log</string>
	<key>StandardOutPath</key>
	<string>REPLACE_WITH_DATA_DIR/logs/stdout.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>8</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>WorkingDirectory</key>
	<string>REPLACE_WITH_PROJECT_DIR</string>
</dict>
</plist>
